# com.local.writing-tool.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.writing-tool</string>

    <key>ProgramArguments</key>
    <array>
        <!-- UPDATE THIS PATH to match your setup -->
        <string>/Users/YOUR_USERNAME/writing-tool/start.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/writing-tool.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/writing-tool.err</string>
</dict>
</plist>
